# haveged-stop/DEBIAN/postinst
#!/bin/sh
set -e

     
    systemctl daemon-reload
    systemctl enable haveged-stop



exit 0
